# super_agent_bridge/rootfs/etc/s6-overlay/s6-rc.d/super-agent-bridge/run
#!/usr/bin/with-contenv bashio
exec python3 /run.py
